# v2/common/runtime/s6/variants/swoole/services.d/swoole/run
#!/bin/sh
set -e

: "${SWOOLE_HOST:=0.0.0.0}"
: "${SWOOLE_PORT:=80}"
: "${SWOOLE_DOCUMENT_ROOT:=/var/www/html/public}"
: "${SWOOLE_WORKER_NUM:=auto}"
: "${SWOOLE_STATIC:=1}"
: "${ENABLE_OCTANE:=false}"
: "${SWOOLE_WATCH:=0}"
: "${DISABLE_XDEBUG_IN_SWOOLE:=1}"

# Optional: ensure Xdebug is off for Swoole workers unless you explicitly enable it
if [ "$DISABLE_XDEBUG_IN_SWOOLE" = "1" ]; then
  for f in /usr/local/etc/php/conf.d/*xdebug*.ini; do
    [ -f "$f" ] && mv "$f" "${f}.disabled" || true
  done
fi

if [ "$ENABLE_OCTANE" = "true" ] && [ -f /var/www/html/artisan ]; then
  echo "----> Starting Laravel Octane (Swoole)..."
  exec /common/s6-setuidgid www-data \
    php -d opcache.enable_cli=0 /var/www/html/artisan octane:start \
      --server=swoole --host="$SWOOLE_HOST" --port="$SWOOLE_PORT" \
      --workers="${SWOOLE_WORKER_NUM:-auto}" ${OCTANE_ARGS:-}
fi

if php -m | grep -qiE '^openswoole$|^swoole$'; then
  echo "----> Starting Swoole HTTP server..."
  if [ "$SWOOLE_WATCH" = "1" ]; then
    # Dev: force OPcache off at the CLI level so edits show instantly
    exec /common/s6-setuidgid www-data \
      php -d opcache.enable_cli=0 /usr/local/bin/mani-swoole-server
  else
    exec /common/s6-setuidgid www-data \
      php /usr/local/bin/mani-swoole-server
  fi
fi

echo "[swoole] extension not loaded; sleeping"
exec tail -f /dev/null
