# v2/common/runtime/s6/variants/nginx/services.d/queue-worker/run
#!/command/with-contenv /bin/bash
set -euo pipefail
if [[ ! -f /var/www/html/artisan ]]; then
  echo "[queue-worker] artisan not found; sleeping"
  exec tail -f /dev/null
fi
echo "----> Starting Laravel queue worker..."
# exec /usr/bin/s6-setuidgid www-data \
exec /common/s6-setuidgid www-data \
  php /var/www/html/artisan queue:work --verbose --tries=3 --timeout=90
