# v1/common/supervisord.apache.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autostart=true
autorestart=true
stdout_logfile=/var/log/apache2/apache2.stdout.log
stderr_logfile=/var/log/apache2/apache2.stderr.log
exitcodes=0


; Optional: program for Laravel Horizon or other queue workers
; Uncomment and configure if you need queue workers managed by supervisord.
; [program:queue-worker]
; process_name=%(program_name)s_%(process_num)02d
; command=php /var/www/html/artisan queue:work --sleep=3 --tries=3
; autostart=false ; Set to true if you want the worker to start automatically
; autorestart=true
; user=www-data ; Run the worker as the web server user
; numprocs=1 ; Number of worker processes
; redirect_stderr=true
; stdout_logfile=/var/www/html/storage/logs/worker.log ;
