# v2/common/runtime/s6/variants/apache/services.d/apache/run
#!/bin/sh
set -e

. /etc/apache2/envvars

echo "----> Starting Apache2 service..."
exec /usr/sbin/apache2ctl -D FOREGROUND
