# v2/common/runtime/s6/variants/nginx-swoole/services.d/swoole/run
#!/bin/sh
set -e
: "${SWOOLE_HOST:=127.0.0.1}"
: "${SWOOLE_PORT:=9501}"                 # nginx proxies to this
: "${SWOOLE_DOCUMENT_ROOT:=/var/www/html/public}"
: "${SWOOLE_WORKER_NUM:=auto}"
: "${SWOOLE_STATIC:=1}"
: "${ENABLE_OCTANE:=false}"

if [ "$ENABLE_OCTANE" = "true" ] && [ -f /var/www/html/artisan ]; then
  echo "----> Starting Laravel Octane (Swoole) on ${SWOOLE_HOST}:${SWOOLE_PORT}..."
  exec /command/s6-setuidgid www-data \
    php /var/www/html/artisan octane:start \
      --server=swoole --host="$SWOOLE_HOST" --port="$SWOOLE_PORT" \
      ${OCTANE_ARGS:-}
fi

if php -m | grep -qiE '^openswoole$|^swoole$'; then
  echo "----> Starting Swoole HTTP (nginx front) on ${SWOOLE_HOST}:${SWOOLE_PORT}"
  exec /command/s6-setuidgid www-data /usr/local/bin/mani-swoole-server
fi

echo "[swoole] extension not loaded; sleeping"
exec tail -f /dev/null
