# v2/tests/fixtures/overlays/hello-svc/services.d/hello/run
#!/command/with-contenv /bin/bash
set -euo pipefail
echo "[hello] service up"
# prove we’re alive by touching a file repeatedly
while true; do
  date -u > /tmp/hello.probe
  sleep 5
done
