# v2.1/examples/.docker-snippets/overlays/helloworld-svc/services.d/hello-wolrd/finish
#!/usr/bin/env bash

echo "[hello-world] Custom service is shutting down."
exit 0
