# v2/common/runtime/s6/variants/nginx-swoole/services.d/nginx/run
#!/bin/sh
set -e
echo "----> Starting Nginx..."
exec /usr/sbin/nginx -g 'daemon off;'
