# v2/common/runtime/s6/variants/nginx/services.d/nginx/run
#!/bin/sh
set -e
echo "----> Starting Nginx service..."
# ensure pid dir
mkdir -p /run/nginx
exec nginx -g 'daemon off;'
